# scripts/com.hagent.connect-remote-disk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hagent.connect-remote-disk</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nguyenhat/HAgent/scripts/connect-remote-disk.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Library/Preferences/SystemConfiguration/com.apple.identification.plist</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nguyenhat/HAgent/logs/connect-remote-disk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nguyenhat/HAgent/logs/connect-remote-disk.log</string>
</dict>
</plist>
