# scripts/com.hagent.auto-shutdown.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hagent.auto-shutdown-remote</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nguyenhat/HAgent/scripts/auto-shutdown-remote.py</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nguyenhat/HAgent/logs/auto-shutdown-remote.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nguyenhat/HAgent/logs/auto-shutdown-remote.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SSH_PASSWORD</key>
        <string>NgocNhi@1811</string>
    </dict>
</dict>
</plist>
